# bishop-meta/launchd/com.bishopbot.local-worker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.bishopbot.local-worker</string>

    <key>WorkingDirectory</key>
    <string>/Users/matthewbishop/BishopBot</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/matthewbishop/BishopBot/bishop-meta/local_worker_daemon.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/matthewbishop/BishopBot/logs/launchd-local-worker.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/matthewbishop/BishopBot/logs/launchd-local-worker.err.log</string>
  </dict>
</plist>
